# infra/systemd/jarv1s-backup.service
# NOTE: `User=REPLACE_WITH_INSTALL_USER` and the `~/Jarv1s` paths are PLACEHOLDERS
# the installing operator MUST edit to match their system (user, repo path, and
# EnvironmentFile path) at install time.
[Unit]
Description=Jarv1s scheduled backup (DB + vault)
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=REPLACE_WITH_INSTALL_USER
WorkingDirectory=~/Jarv1s
EnvironmentFile=-~/Jarv1s/infra/backup.env
ExecStart=~/Jarv1s/scripts/backup-full.sh
StandardOutput=journal
StandardError=journal
