# infra/systemd/jarv1s-stack.service
# NOTE: `User=REPLACE_WITH_INSTALL_USER` and the `~/Jarv1s` paths follow the EXISTING repo
# convention (`infra/systemd/jarv1s-backup.service` hardcodes the same) — this is a
# single-operator household deploy (ADR 0007), not a portable distro package. An
# operator on a different path/user edits these two lines (and the EnvironmentFile
# path) at install time. The committed file IS the canonical install for this host.
[Unit]
Description=Jarv1s production stack (docker compose)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=REPLACE_WITH_INSTALL_USER
WorkingDirectory=~/Jarv1s
# Operator env file: provides JARVIS_IMAGE_TAG (REQUIRED — the prod Compose errors
# without it), host UID/GID, socket paths, and secrets for `docker compose ${...}`
# interpolation. We pass it to docker via `--env-file` (NOT systemd `EnvironmentFile=`)
# for two reasons: (1) it exactly matches the documented manual deploy command
# (`docker compose --env-file <file> ... up -d`), and (2) systemd's EnvironmentFile
# parser and docker's --env-file parser treat `$` differently — a generated secret
# (e.g. POSTGRES_PASSWORD / BETTER_AUTH_SECRET) containing `$` would need `$$` for
# systemd but stays literal for docker, silently corrupting auth/crypto on reboot.
# Letting docker read the file directly avoids that whole class of `$`-mangling.
# The operator MUST create this file (with JARVIS_IMAGE_TAG) before `systemctl enable`;
# if it is missing, the Compose `${...:?}` guards fail the start LOUDLY.
ExecStart=/usr/bin/docker compose --env-file ~/Jarv1s/infra/env.production.local -f infra/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose --env-file ~/Jarv1s/infra/env.production.local -f infra/docker-compose.prod.yml down
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
